# systemd/podripper-radiot.timer
[Timer]
OnCalendar=Sat 22:50 Europe/Moscow
AccuracySec=1m

[Install]
WantedBy=timers.target
