# systemd/podripper@radiot.timer
[Timer]
OnCalendar=Sat 19:50 UTC
AccuracySec=1m

[Install]
WantedBy=timers.target
