# systemd/podripper@.service
[Unit]
Description=Podripper %i
After=network.target

[Service]
ExecStart=/usr/bin/ripper-exe run %i
WorkingDirectory=/var/lib/podripper
User=user

[Install]
WantedBy=multi-user.target
