# systemd/podripper-rcmp.timer
[Timer]
OnCalendar=Thu 16:30 America/Toronto
AccuracySec=1m

[Install]
WantedBy=timers.target
